# zfs-backup/etc/zfs-backup.timer
[Unit]
Description=timer to backup zfs pool 

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target