# zfs-backup/etc/zfs-backup.service
[Unit]
Description=backup zfs pool onto remote

[Service]
ExecStart=/usr/local/bin/zfs-backup.sh

[Install]
WantedBy=multi-user.target
