# zfs-backup/etc/zfs-backup.timer
[Unit]
Description=timer to backup zfs pool 
Requires=zfs.target
After=zfs.target

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
