# google-drive/google-drive-sync.timer
[Unit]
Description=Timer for google-drive sync

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
