# zfs-auto-scrub/zfs-scrub.service
[Unit]
Description=scrub zpool0

[Service]
Type=oneshot
ExecStart=zpool scrub zpool0
