# zfs-auto-scrub/zfs-scrub.timer
[Unit]
Description=weekly timer for zfs scrub

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=multi-user.target
