# zfs-auto-scrub/etc/zfs-scrub@.service
[Unit]
Description=scrub zfs pool %i

[Service]
Type=oneshot
ExecStart=zpool scrub %i
