# zfs-auto-scrub/etc/zfs-scrub@.timer
[Unit]
Description=weekly timer for zfs for %i

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
